# fox.service
[Unit]
Description=Ferrumox LLM inference server (%i)
After=network.target

[Service]
ExecStart=/usr/local/bin/fox serve --model-path %i
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
